# xdg-config/autostart/openvpn3-applet.desktop.in
[Desktop Entry]
Exec=sh -c '@execPath@'
Icon=
Name=openvpn3-applet
Path=
Terminal=False
Type=Application